# config/upstart/wlm-sneakers.conf
start on starting wlm
stop on stopping wlm
